# node_exporter.service
[Unit]
Description=Node exporter for VictoriaMetrics
[Service]
ExecStart=/opt/node_exporter/node_exporter
[Install]
WantedBy=multi-user.target