# node_exporter.service
[Unit]
Description=Node exporter for VictoriaMetrics
[Service]
ExecStart=/usr/local/bin/node_exporter
[Install]
WantedBy=multi-user.targetc