# templates/node_exporter.service.j2
[Unit]
Description=Node exporter for VictoriaMetrics
[Service]
ExecStart={{node_exporter_dir}}/node_exporter --web.listen-address= {{ node_exporter_web_address }}
[Install]
WantedBy=multi-user.target